# systemd/hal-voice-assistant.service
[Unit]
Description=HAL 9000 Voice Assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/hal-voice-assistant
EnvironmentFile=/opt/hal-voice-assistant/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/hal-voice-assistant/.pixi/envs/default/bin/hal-voice-assistant
Restart=on-failure
RestartSec=5
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
